# launchagent/com.claude-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-watchdog</string>

    <!-- Network + poller-health watchdog. Runs the connectivity loop and,
         on a stale/absent poller heartbeat, force-restarts the poller job
         via `launchctl kickstart -k com.claude-agent`. -->
    <key>ProgramArguments</key>
    <array>
        <string>__BIN_DIR__/watch_network.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Respawn the watchdog only on a CRASH (non-zero exit), NOT after a clean
         shutdown. watch_network.sh traps SIGTERM and exits 0, so `cta stop`
         (launchctl unload → SIGTERM) leaves the watchdog DOWN instead of being
         immortally respawned. A genuine crash (non-zero exit) still relaunches
         it. This fixes the immortal-watchdog bug where an unconditional
         KeepAlive re-kicked the poller back to life after every intentional
         stop, and the orphaned caffeinate kept the Mac awake. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- PATH must include:
             - ~/.bun/bin — bun
             - /opt/homebrew/bin — Apple Silicon Homebrew
             - /usr/local/bin — Intel Homebrew
             so the same plist works on both architectures. -->
        <key>PATH</key>
        <string>__HOME__/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <!-- Watchdog connectivity logs join the poller's agent.log so Pager's
         "Show recent activity" tail surfaces both streams. -->
    <key>StandardOutPath</key>
    <string>__STATE_DIR__/agent.log</string>

    <key>StandardErrorPath</key>
    <string>__STATE_DIR__/agent.err.log</string>
</dict>
</plist>
