# launchagent/com.claude-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BIN_DIR__/start_agents.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Respawn only on unsuccessful exit (crash). A clean `cta stop`
         (launchctl unload) is treated as intentional and the agent stays
         dead until the next login or explicit `cta start`. Matches the
         Pager LaunchAgent's policy. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- PATH must include:
             - ~/.bun/bin — bun (Telegram plugin spawns `bun run …`; if missing,
               the channel silently fails to start)
             - /opt/homebrew/bin — Apple Silicon Homebrew (tmux, etc.)
             - /usr/local/bin — Intel Homebrew (tmux, etc.)
             so the same plist works on both architectures. -->
        <key>PATH</key>
        <string>__HOME__/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__STATE_DIR__/agent.log</string>

    <key>StandardErrorPath</key>
    <string>__STATE_DIR__/agent.err.log</string>
</dict>
</plist>
