# launchagent/com.claude-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-watchdog</string>

    <!-- Network + poller-health watchdog. Runs the connectivity loop and,
         on a stale/absent poller heartbeat, force-restarts the poller job
         via `launchctl kickstart -k com.claude-agent`. -->
    <key>ProgramArguments</key>
    <array>
        <string>__BIN_DIR__/watch_network.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Always keep the watchdog alive. Unlike the poller (which stays down
         after a clean `cta stop`), the watchdog is pure supervision with no
         conversational state, so an unconditional KeepAlive is correct: if it
         dies for any reason, relaunch it. A clean `cta stop` boots it out
         explicitly so this never fights an intentional shutdown. -->
    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- PATH must include:
             - ~/.bun/bin — bun
             - /opt/homebrew/bin — Apple Silicon Homebrew
             - /usr/local/bin — Intel Homebrew
             so the same plist works on both architectures. -->
        <key>PATH</key>
        <string>__HOME__/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <!-- Watchdog connectivity logs join the poller's agent.log so Pager's
         "Show recent activity" tail surfaces both streams. -->
    <key>StandardOutPath</key>
    <string>__STATE_DIR__/agent.log</string>

    <key>StandardErrorPath</key>
    <string>__STATE_DIR__/agent.err.log</string>
</dict>
</plist>
